# makemkv-wrapper/rootfs/etc/services.d/makemkv/finish
#!/usr/bin/execlineb -S0
##
## /etc/services.d/makemkv/finish
## Called by s6 when the makemkv service exits.
## Brings the whole container down so the HA supervisor can restart it cleanly.
##
## ⚠️  s6-overlay v3 ONLY — uses /run/s6/basedir/bin/halt
##     (NOT /run/s6/basedir/bin/s6-svscanctl or /run/s6/services which are v2)
##     Official HA docs: https://developers.home-assistant.io/blog/2022/05/12/s6-overlay-base-images/
##
/run/s6/basedir/bin/halt
