# makemkv-wrapper/rootfs/etc/services.d/makemkv/finish
#!/usr/bin/execlineb -S0
##
## /etc/services.d/makemkv/finish
## Called when the service exits. Brings the whole container down so
## Home Assistant can restart it cleanly.
##
## s6-overlay v3: scan dir is /run/service; binary is at /command/s6-svscanctl.
## The old /run/s6/basedir path only exists in s6-overlay v2.
##
/command/s6-svscanctl -t /run/service
