# makemkv-wrapper/rootfs/etc/services.d/makemkv/run
#!/usr/bin/execlineb -P
##
## /etc/services.d/makemkv/run
## s6 service definition — starts the jlesage startapp process
## Ported by matsob0123
##
with-contenv
/startapp.sh
